# 2023/DarthBob112112/chall2/crontab
*/1 * * * * /bin/bash /home/abhi/dockermonitor.sh

